# etc/cron.d
# Creates zfs snapshot time creation log file
#
#min    hours   DayOfMonth      Month   DayOfWeek user     command
# Incr backup mon-tue
0      *      *               *      *        zabbix /etc/zabbix/scripts/zabbix-zfs-snapshot
